# deploy/systemd/lms-backup-sync.timer
[Unit]
Description=Sync LMS Portal backups to off-server storage daily

[Timer]
OnCalendar=*-*-* 03:50:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
