# deploy/systemd/lms-backup-sync.service
[Unit]
Description=LMS Portal off-server backup sync
Wants=network-online.target
After=network-online.target lms-backup-verify.service

[Service]
Type=oneshot
WorkingDirectory=/var/www/lms-portal
Environment=ENV_FILE=/var/www/lms-portal/Backend/.env.production
ExecStart=/bin/sh /var/www/lms-portal/Backend/scripts/backup/sync-backups.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
